# supervisord-newrelic.conf
[program:newrelic]
command=/usr/sbin/nrsysmond -c /etc/newrelic/nrsysmond.cfg -l /dev/stdout -f
numprocs=1
autostart=true
autorestart=true
